# services/launchd/com.obsidian-tools.indexer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.obsidian-tools.indexer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/obsidian-tools/.venv/bin/python</string>
        <string>src/index_vault.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/obsidian-tools</string>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/obsidian-tools-indexer.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/obsidian-tools-indexer.log</string>
</dict>
</plist>
