# services/systemd/obsidian-tools-indexer-scheduler.timer
[Unit]
Description=Run Obsidian Tools Vault Indexer hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
