# services/systemd/obsidian-tools-indexer.service
[Unit]
Description=Obsidian Tools Vault Indexer

[Service]
Type=oneshot
WorkingDirectory=__PROJECT_DIR__
ExecStart=__VENV_PYTHON__ src/index_vault.py
