# services/systemd/obsidian-tools-api.service
[Unit]
Description=Obsidian Tools API Server
After=network.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
ExecStart=__VENV_PYTHON__ src/api_server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
